# deploy/systemd/minecraft@.service
[Unit]
Description=Minecraft server %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=minecraft
Group=minecraft
WorkingDirectory=/srv/minecraft/%i
EnvironmentFile=-/etc/minecraft/%i.env
ExecStart=/usr/local/sbin/start-minecraft-server %i
SuccessExitStatus=0 143
Restart=on-failure
RestartSec=10
TimeoutStopSec=90
UMask=0007

[Install]
WantedBy=multi-user.target
